# etc/init/pmp_client.conf
description     "pmp_client daemon"

start on runlevel [2345]
stop on runlevel [016]

#expect fork
expect daemon
respawn
console none
exec pmp_client
